# docker-neo-csharp-nodes/scripts/taillogs.cron
* * * * * /opt/taillogs.sh > /dev/null

